# root/etc/init.d/todo-101
#!/sbin/openrc-run

name="todo"
description="A simple todo list manager"
command="/usr/bin/node"
directory="/app"
command_args="src/index.js"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/var/log/${RC_SVCNAME}.log"
error_log="${output_log}"

depend() {
  need net
}
